# target/linux/ramips/base-files-ramips/etc/inittab
# Default entries for busybox init.
::ctrlaltdel:/sbin/reboot
#::shutdown:/bin/umount -a -r
::shutdown:/etc/init.d/rcD
::restart:/sbin/init
#ttyS1::askfirst:-/bin/sh
::sysinit:/etc/init.d/rcS
#::respawn:-/bin/sh
::askfirst:/sbin/getty -L ttyS1 57600 vt100
